# Ejemplos_practicos/Targets/modo1_con_red.target
[Unit]
Description=Modo 1 con red
# Documentation=man:systemd.special(7)
Requires=rescue.target
Wants=NetworkManager.service
After=rescue.target
AllowIsolate=yes
